# packaging/debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop service
        if systemctl is-active --quiet tomo 2>/dev/null; then
            systemctl stop tomo || true
        fi

        # Disable service (but don't remove unit file yet)
        systemctl disable tomo 2>/dev/null || true
        ;;
    failed-upgrade)
        ;;
    *)
        echo "prerm called with unknown argument: $1" >&2
        exit 1
        ;;
esac

exit 0
